# cortexm_uclinux-master/filesystem/rootfs/local/inittab
::sysinit:/etc/rc
ttyS0::wait:mdev -s
ttyS0::once:/bin/hush -i
ttyS0::wait:/bin/mount_jffs2.sh
ttyS0::once:udhcpc
ttyS0::once:udhcpc
ttyS0::once:udhcpc
ttyS0::once:udhcpc
ttyS0::once:udhcpc
ttyS0::wait:cp /mnt/JointWarn ./
ttyS0::once:./JointWarn
ttyS0::once:echo 100 > sys/class/backlight/pwm-backlight.0/brightness
