# cortexm_uclinux-master/filesystem/rootfs/local/inittab
::sysinit:/etc/rc
ttyS0::wait:mdev -s
ttyS0::once:/bin/hush -i
ttyS0::once:/bin/fbtest
ttyS0::once:/bin/mount_jffs2.sh
