# cortexm_uclinux-master/filesystem/rootfs/local/inittab
::sysinit:/etc/rc
ttyS0::respawn:/bin/hush -i
ttyS0::respawn:mdev -s
